# deployment/supervisor/uwsgi.conf
[program:uwsgi]
user=root
command = /usr/local/bin/uwsgi --ini /duole_mini/deployment/uwsgi/uwsgi.ini
autostart=true
autorestart=true
startsecs=10